# deploy/pi-watcher.service
[Unit]
Description=pi-watcher network monitor
# time-sync.target matters on RTC-less boards (every Raspberry Pi): starting
# before NTP sync stamps samples and durations with a bogus clock. For it to
# actually wait, enable the sync unit once:  sudo systemctl enable systemd-time-wait-sync
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target

[Service]
ExecStart=/usr/local/bin/pi-watcher --config /etc/pi-watcher/config.toml
Restart=always
RestartSec=5

# Raw ICMP sockets without running as root.
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW

DynamicUser=yes
StateDirectory=pi-watcher
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
